# builder/overlay-rootfs/etc/s6/services/docker-cmd/finish
#!/bin/execlineb -S0
redirfd -w 1 /var/run/s6/env-stage3/S6_CMD_EXITED s6-echo -n -- "${1}"
