# builder/overlay-rootfs/etc/s6/services/docker-cmd/finish
#!/bin/execlineb -S0
foreground { redirfd -w 1 /var/run/s6/env-stage3/S6_STAGE2_EXITED s6-echo -n -- "${1}" }
foreground { s6-echo -- "docker CMD exited ${1}" }
s6-svscanctl -t /var/run/s6/services
