# overlay-rootfs/etc/s6/fdholderd/run
#!/usr/bin/execlineb
s6-notifywhenup s6-fdholder-daemon -1 -i /etc/s6/fdholderd/rules /etc/s6/fdholderd/socket