# builder/overlay-rootfs/etc/s6/service/s6-svscan-log/run
#!/usr/bin/execlineb -P
redirfd -rnb 0 /var/run/s6/uncaught-logs-fifo
s6-envuidgid -D 32768:32768 nobody
s6-applyuidgid -U
exec -c
s6-log -bp -- t /var/log/s6-uncaught-logs
