# overlay-rootfs/etc/s6/.s6-svscan/finish
#!/usr/bin/execlineb
/etc/s6/.s6-init/init-stage3
