# layout/syslogd-overlay/etc/s6-overlay/s6-rc.d/syslogd-log/run
#!@SHEBANGDIR@/execlineb -P

s6-envuidgid -D 32761:32761: -- sysllog
s6-applyuidgid -U --
backtick -D "n20 s1000000" line { printcontenv S6_LOGGING_SCRIPT }
multisubstitute
{
  importas -usC T line
  define dir /var/log/syslogd
}
exec -c
s6-log -d3 --

 - +^auth\\. +^authpriv\\.  $T ${dir}/secure
 - +^cron\\.  $T ${dir}/cron
 - +^daemon\\.  $T ${dir}/daemon
 - +^[[:alnum:]]*\\.debug:  $T ${dir}/debug
 - +^[[:alnum:]]*\\.err: +^[[:alnum:]]*\\.error: +^[[:alnum:]]*\\.emerg: +^[[:alnum:]]*\\.alert: +^[[:alnum:]]*\\.crit:  $T ${dir}/errors
 - +^kern\\.  $T ${dir}/kernel
 - +mail\\.  $T ${dir}/mail
 - +user\\.  $T ${dir}/user
 - +^[[:alnum:]]*\\.info: +^[[:alnum:]]*\\.notice: +^[[:alnum:]]*\\.warn: -^auth\\. -^authpriv\\. -^cron\\. -daemon\\. -mail\\.  $T ${dir}/messages
 + -^auth\\. -^authpriv\\.  $T ${dir}/everything
