# builder/overlay-rootfs/etc/s6/services/s6-svscan-log/run
#!/bin/execlineb -P
logutil-service -f /var/run/s6/uncaught-logs-fifo /var/log/s6-uncaught-logs
