# builder/overlay-rootfs/etc/s6/services/s6-svscan-log/run
#!/usr/bin/execlineb -P
logutil-service /var/run/s6/uncaught-logs-fifo /var/log/s6-uncaught-logs
