# layout/syslogd-overlay/etc/s6-overlay/s6-rc.d/syslogd/run
#!@SHEBANGDIR@/execlineb -P

s6-envuidgid -D 32760:32760: -- syslog
s6-socklog -d3 -U -t3000
